# opensips/units/etc/service/opensips/run
#!/bin/bash
set -m

# Auto-scale based on target CC
TARGET_CC="${TARGET_CC:-2000}"

if [ "$TARGET_CC" -le 2000 ]; then
    SHM_SIZE="${OPENSIPS_SHM_SIZE:-768}"
    PKG_SIZE="${OPENSIPS_PKG_SIZE:-64}"
elif [ "$TARGET_CC" -le 3500 ]; then
    SHM_SIZE="${OPENSIPS_SHM_SIZE:-1024}"
    PKG_SIZE="${OPENSIPS_PKG_SIZE:-64}"
elif [ "$TARGET_CC" -le 6000 ]; then
    SHM_SIZE="${OPENSIPS_SHM_SIZE:-2048}"
    PKG_SIZE="${OPENSIPS_PKG_SIZE:-128}"
else
    SHM_SIZE="${OPENSIPS_SHM_SIZE:-4096}"
    PKG_SIZE="${OPENSIPS_PKG_SIZE:-128}"
fi

echo "Starting OpenSIPS for ~${TARGET_CC} CC (SHM=${SHM_SIZE}MB, PKG=${PKG_SIZE}MB)"
exec opensips -f ${OPENSIPS_CFG} -F -m "$SHM_SIZE" -M "$PKG_SIZE"