# opensips/units/etc/service/opensips/run
#!/bin/bash
set -m

/usr/sbin/opensips -f ${OPENSIPS_CFG} -F