# rtpproxy/units/etc/service/rtpproxy/run
#!/bin/bash

set -e
set -o pipefail

init_vars() {
  RTPPROXY_LISTEN_HOST=${RTPPROXY_LISTEN_HOST:-"0.0.0.0"}
  RTPPROXY_CTRL_SOCKET=${RTPPROXY_CTRL_SOCKET:-"udp:*:7890"}
  RTPRPOXY_NOFILE_LIMIT=${RTPRPOXY_NOFILE_LIMIT:-500000}
  RTPPROXY_MIN_PORT=${RTPPROXY_MIN_PORT:-10000}
  RTPPROXY_MAX_PORT=${RTPPROXY_MAX_PORT:-65000}
  RTPPROXY_LOG_LEVEL=${RTPPROXY_LOG_LEVEL:-"INFO:LOG_LOCAL5"}
  local cmd="$RTPPROXY_BIN -RFf -l $RTPPROXY_LISTEN_HOST -s $RTPPROXY_CTRL_SOCKET -L $RTPRPOXY_NOFILE_LIMIT -m $RTPPROXY_MIN_PORT -M $RTPPROXY_MAX_PORT -r /var/lib/rtpproxy-recording -S /var/spool/rtpproxy"
  RTPPROXY_CMD=${RTPPROXY_CMD:-"$cmd"}
}

main() {
  if [ -z $RTPPROXY_PARAM ]; then
    init_vars
  else
    RTPPROXY_CMD="$RTPPROXY_BIN $RTPPROXY_PARAM"
  fi
  exec $RTPPROXY_CMD
}

main